# net-vpn/ivpn/files/ivpn.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="IVPN-Service Daemon"
pidfile="${RC_PREFIX}/run/ivpn.pid"

command="/usr/bin/ivpn-service"
command_user="root:root"

reload() {
	ebegin "Reloading ${RC_SVCNAME}"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}
