# net-vpn/ivpn/files/ivpn.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="IVPN-Service Daemon"

command="/usr/bin/ivpn-service"
command_background=true
command_user="root:root"

supervisor="supervise-daemon"
supervise_daemon_args=" -d /"

depend() {
    need localmount net
    use dns
    after bootmisc
    provide vpn
}
