# app-misc/input-remapper/files/input-remapper.initd
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Input-Remapper"
description="Daemon for input-remapper service"
pidfile="/var/run/input-remapper.pid"
command="/usr/bin/input-remapper-service"
command_background="yes"

depend() {
	need dbus
}
