# sys-apps/touchegg/files/touchegg.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="${RC_PREFIX}/run/touchegg.pid"
name="Touchégg Daemon"
command="${RC_PREFIX}/usr/bin/touchegg"
command_args="--daemon"
command_background="true"
command_user="root:input"
